# deploy/gitrepo/crontab

SHELL=/bin/bash
MAILTO=""

0 */2 * * * /home/git/autopush.sh
