# deploy/ddns/crontab

SHELL=/bin/bash
MAILTO=""

0 * * * * /data/ddns/ddns.sh
