# deploy/ddns/crontab

SHELL=/bin/bash
MAILTO=""

*/15 * * * * /data/ddns/ddns.sh
